# backend/deploy/vidyatrack.service
[Unit]
Description=VidyaTrack FastAPI backend
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/app/demo-poc/school/backend
Environment="PATH=/home/ubuntu/app/demo-poc/school/backend/venv/bin"
ExecStart=/home/ubuntu/app/demo-poc/school/backend/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 2
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
